# entrypoint/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── Hermes gateway (Telegram/Discord/WhatsApp messaging) ─────────────────────
; Messaging gateway only — NOT the chat API. Runs for platform integrations.
[program:gateway]
command=/opt/hermes/.venv/bin/hermes gateway run --accept-hooks
; CWD is /opt/data/workspace so any tool subprocess that writes a relative
; filename (e.g. camoufox screenshot("foo.png")) lands in the bind-mounted
; workspace dir — visible in the Files tab and miniserve, persisted to host.
directory=/opt/data/workspace
user=workspace
environment=HOME=/home/workspace,HERMES_HOME=/opt/data,GOOGLE_WORKSPACE_CLI_CREDENTIALS_FILE=/home/workspace/.config/gws/credentials.json,PATH="/opt/hermes/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Hermes Dashboard (port 9119) ─────────────────────────────────────────────
; Hermes web dashboard — provides /api/skills, /api/sessions, /api/config on
; loopback:9119. hermes-workspace probes this to unlock Skills, Sessions, and
; Agent Settings. Must start before workspace (priority 15).
[program:dashboard]
command=/opt/hermes/.venv/bin/hermes dashboard --no-open --host 127.0.0.1 --port 9119
; Same rationale as the gateway — see comment above.
directory=/opt/data/workspace
user=workspace
environment=HOME=/home/workspace,HERMES_HOME=/opt/data,GOOGLE_WORKSPACE_CLI_CREDENTIALS_FILE=/home/workspace/.config/gws/credentials.json,PATH="/opt/hermes/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
autostart=true
autorestart=true
priority=15
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Hermes Workspace UI (port 8787) ───────────────────────────────────────────
; Node.js server for the hermes-workspace React UI.
; HERMES_API_URL → gateway API server (port 8642); HERMES_DASHBOARD_URL → dashboard (port 9119).
[program:workspace]
command=sh -c '/opt/hermes/scripts/wait-for-hermes-gateway.sh && exec node /app/server-entry.js'
directory=/app
user=workspace
environment=HOME=/home/workspace,PORT=8787,HOST=0.0.0.0,HERMES_API_URL="http://127.0.0.1:8642",HERMES_DASHBOARD_URL="http://127.0.0.1:9119",HERMES_HOME="/opt/data",NODE_ENV=production,COOKIE_SECURE=0,CLAUDE_PASSWORD="%(ENV_CLAUDE_PASSWORD)s",HERMES_API_TOKEN="%(ENV_HERMES_API_TOKEN)s",PYTHONPATH="/opt/hermes/.venv/lib/python3.11/site-packages",GOOGLE_WORKSPACE_CLI_CREDENTIALS_FILE=/home/workspace/.config/gws/credentials.json,GITHUB_TOKEN="%(ENV_GITHUB_TOKEN)s",PATH="/opt/hermes/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── File explorer (port 8800) ─────────────────────────────────────────────────
; miniserve serves /opt/data/workspace at /__explorer so users can download
; exported files (CSV, JSON, etc.) from the chat UI.
[program:explorer]
command=miniserve --port 8800 --route-prefix /__explorer --no-symlinks --hide-version-footer /opt/data/workspace
user=workspace
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Browser terminal — ttyd (port 7681) ─────────────────────────────────────
; Previously started on-demand by ensure_ttyd_running() during page render.
; That created a race: hitting /__terminal before the page render landed (e.g.
; bookmark, refresh after container restart) → port 7681 not yet listening →
; Traefik returns 502 → Cloudflare shows "Bad gateway". Supervising it as
; always-on eliminates the race for ~3MB of resident memory. ensure_ttyd_running
; is now a fast no-op (pgrep finds the process immediately).
[program:terminal]
command=ttyd --port 7681 --base-path /__terminal -W bash
user=workspace
environment=HOME=/home/workspace,PATH="/opt/hermes/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",TERM=xterm-256color
autostart=true
autorestart=true
priority=25
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── AutoBricks MCP exec-bridge (port 9876) ──────────────────────────────────
; Exposes run_command / read_file / write_file over MCP at https://<bot>/__mcp
; so the user's local Claude Code (or Cursor / Zed) can drive shell + file ops
; on this bot. Auth via ABAI_MCP_TOKEN (loaded from /opt/data/.env at request
; time, so rotations take effect without restarting this process).
; See docker/abai-mcp-server.py for the implementation.
[program:abai-mcp]
command=/opt/hermes/.venv/bin/python /opt/abai-mcp-server.py
user=workspace
environment=HOME=/home/workspace,PATH="/opt/hermes/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",ABAI_MCP_PORT=9876
autostart=true
autorestart=true
priority=35
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
